# infrastructure/automation/resume-sync.service
[Unit]
Description=Resume Platform Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/jclee/dev/apps/resume/job-automation-mcp
ExecStart=/usr/bin/node scripts/auto-all.js --sync
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
